# internal/controller/crontab/crontab
9 11 1,18 4 1 tar -zcf /var/backups/home.tgz /home/
@midnight /home/maverick/bin/cleanup-logs.sh
    11 11 1,10-16,18 4 2-5  /home/ubuntu/command.sh > /tmp/output 2>&1